# day-2-linux-advanced/monitor.service
[Unit]
Description=System Monitor Service
After=network.target

[Service]
Type=simple
ExecStart=/opt/monitor/monitor.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target